# scripts/systemd-calibre-compose.service
[Unit]
Description=Calibre Web (Docker Compose)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/ebooks_lv
ExecStart=/usr/bin/docker compose -f compose.yml -f compose.droplet.yml up -d
ExecStop=/usr/bin/docker compose -f compose.yml -f compose.droplet.yml down
RemainAfterExit=true
TimeoutStartSec=0
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

# Install instructions:
#   sudo cp scripts/systemd-calibre-compose.service /etc/systemd/system/calibre-compose.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable calibre-compose
#   sudo systemctl start calibre-compose